# gosimplefileserver.service
[Unit]
Description=Simple Go File Server

[Service]
Type=simple
ExecStart=go-simple-upload-server -token 123 -port 8888 -upload_limit 100000000 "/home/ansible/projektarbeit"
User=ansible

[Install]
WantedBy=multi-user.target